# docker/php-fpm/supervisor.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work
autostart=false
autorestart=true
numprocs=1
startretries=100000000
startsecs=0
